# rc.d/dhcpleased
#!/bin/ksh

daemon="/usr/local/sbin/dhcpsyncd"

. /etc/rc.d/rc.subr

rc_cmd $1
